# crontab.cron
*/1 * * * * plotmoving\main.py -s D:\\ -d A:\\ -l 1
